# kubernetes/supervisord.conf
[supervisord]
nodaemon=true

[program:etcd]
command=etcd --data-dir=/var/etcd/data